# aur-packages-updater.timer
[Unit]
Description=Check for AUR updates daily

[Timer]
OnCalendar=*-*-* 7:00:00
RandomizedDelaySec=2h

[Install]
WantedBy=timers.target